# cron/crontab
0 * * * * env DISPLAY=:0 /usr/games/xcowsay -t 0 "Coffee Break!"
